# cron/crontab
*/1 * * * * cd /cron; python maragata.py >> /var/log/cron.log 2>&1
