# cron/crontab
30 7 * * * python cron/lunchr_service.py >> /var/log/cron.log 2>&1
